# templates/etc/profile.d/99-terminal_tmout.sh.j2
{{ file_managed_by_ansible }}

#!/bin/bash

{% if ubtu24stig_200060 %}
TMOUT=900
{% endif %}
